# roles/certificates/templates/lego-renewal.timer.j2
[Unit]
Description=Daily LetsEncrypt certificate renewal

[Timer]
OnCalendar=daily
Unit=lego-renewal.service

[Install]
WantedBy=timers.target 